# nginx/boycott.socket
[Unit]
Description=gunicorn socket for boycott.kg

[Socket]
ListenStream=/run/boycott.sock

[Install]
WantedBy=sockets.target